# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=python api_server.py
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/fastapi.err.log
stdout_logfile=/var/log/supervisor/fastapi.out.log
environment=PYTHONUNBUFFERED="1"

[program:streamlit_private]
command=streamlit run streamlit_app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --browser.gatherUsageStats false
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/streamlit_private.err.log
stdout_logfile=/var/log/supervisor/streamlit_private.out.log
environment=PYTHONUNBUFFERED="1"

[program:streamlit_public]
command=streamlit run public_app.py --server.port 8502 --server.address 0.0.0.0 --server.headless true --browser.gatherUsageStats false
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/streamlit_public.err.log
stdout_logfile=/var/log/supervisor/streamlit_public.out.log
environment=PYTHONUNBUFFERED="1"

[program:mcp]
command=python mcp_server.py
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/mcp.err.log
stdout_logfile=/var/log/supervisor/mcp.out.log
environment=PYTHONUNBUFFERED="1"

[program:mcp]
command=python mcp_server.py
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/mcp.err.log
stdout_logfile=/var/log/supervisor/mcp.out.log
environment=PYTHONUNBUFFERED="1"
